# gemini-terminal/rootfs/etc/s6-overlay/s6-rc.d/init-gemini/run
#!/command/with-contenv bashio
# ==============================================================================
# Gemini Terminal - Initialization Service (oneshot)
# Runs once at container startup to prepare the environment
# ==============================================================================

bashio::log.info "Initializing Gemini Terminal..."

# Save addon version for the terminal banner
bashio::addon.version > /data/.addon_version

# ==============================================================================
# Read configuration options
# ==============================================================================

GEMINI_API_KEY=$(bashio::config 'gemini_api_key' || echo "")
MCP_ENABLED=$(bashio::config 'enable_ha_mcp' || echo "true")
SCREENSHOT_ENABLED=$(bashio::config 'screenshot_enabled' || echo "false")
ACCESS_TOKEN=$(bashio::config 'access_token' || echo "")
SMART_CONTEXT=$(bashio::config 'ha_smart_context' || echo "true")
AUTO_LAUNCH=$(bashio::config 'auto_launch_gemini' || echo "true")

bashio::log.info "Configuration: MCP=${MCP_ENABLED} SMART_CONTEXT=${SMART_CONTEXT} AUTO_LAUNCH=${AUTO_LAUNCH}"

# ==============================================================================
# Process user-defined environment variables
# Reads env_vars from add-on config and writes a sourceable file at /data/.env_vars
# ==============================================================================

ENV_VARS_FILE="/data/.env_vars"
rm -f "${ENV_VARS_FILE}"

# Set up Gemini specific environment variables
cat >> "${ENV_VARS_FILE}" << 'EOF'
export HOME="/data"
export XDG_CONFIG_HOME="/data/.config"
export XDG_DATA_HOME="/data/.local/share"
export LANG="en_US.UTF-8"
export LC_ALL="en_US.UTF-8"
export TERM="xterm-256color"
export GEMINI_CONFIG_DIR="/data/.config/gemini"
export GEMINI_HOME="/data"
export GEMINI_TELEMETRY=off
export GEMINI_CLI_NO_RELAUNCH=1
export NODE_OPTIONS="--max-old-space-size=8192"
export UV_THREADPOOL_SIZE=64
export FSWATCH_BACKEND="poll"
EOF

if [ -n "$GEMINI_API_KEY" ] && [ "$GEMINI_API_KEY" != "null" ]; then
    printf "export GOOGLE_API_KEY='%s'\n" "${GEMINI_API_KEY}" >> "${ENV_VARS_FILE}"
    printf "export GEMINI_API_KEY='%s'\n" "${GEMINI_API_KEY}" >> "${ENV_VARS_FILE}"
    bashio::log.info "Gemini API key configured"
fi

ENV_VARS_COUNT=$(bashio::config 'env_vars | length')
if [ "${ENV_VARS_COUNT}" -gt 0 ]; then
    bashio::log.info "Processing ${ENV_VARS_COUNT} user-defined environment variable(s)..."
    for i in $(seq 0 $((ENV_VARS_COUNT - 1))); do
        VAR_NAME=$(bashio::config "env_vars[${i}].name")
        VAR_VALUE=$(bashio::config "env_vars[${i}].value")
        if ! printf '%s' "${VAR_NAME}" | grep -qE '^[A-Za-z_][A-Za-z0-9_]*$'; then
            continue
        fi
        ESCAPED_VALUE=$(printf '%s' "${VAR_VALUE}" | sed "s/'/'\\\\''/g")
        printf "export %s='%s'\n" "${VAR_NAME}" "${ESCAPED_VALUE}" >> "${ENV_VARS_FILE}"
    done
fi

if [ -n "${ACCESS_TOKEN}" ]; then
    printf 'export HA_ACCESS_TOKEN=%q\n' "${ACCESS_TOKEN}" >> "${ENV_VARS_FILE}"
fi

printf 'export HA_URL=http://supervisor/core\n' >> "${ENV_VARS_FILE}"
printf 'export HA_TOKEN=%q\n' "${SUPERVISOR_TOKEN}" >> "${ENV_VARS_FILE}"
printf 'export HA_VERIFY_SSL=false\n' >> "${ENV_VARS_FILE}"

# Z2M configuration
Z2M_URL_CONFIG=$(bashio::config 'z2m_url' || echo "")
Z2M_TOPIC_CONFIG=$(bashio::config 'z2m_mqtt_topic' || echo "")
if [ -n "${Z2M_URL_CONFIG}" ]; then
    printf 'export Z2M_URL=%q\n' "${Z2M_URL_CONFIG}" >> "${ENV_VARS_FILE}"
    if [ -n "${Z2M_TOPIC_CONFIG}" ]; then
        printf 'export Z2M_MQTT_TOPIC=%q\n' "${Z2M_TOPIC_CONFIG}" >> "${ENV_VARS_FILE}"
    fi
else
    Z2M_EXPORTS=$(node /usr/local/bin/discover-z2m.js 2>/dev/null || true)
    if [ -n "${Z2M_EXPORTS}" ]; then
        printf '%s\n' "${Z2M_EXPORTS}" >> "${ENV_VARS_FILE}"
    fi
fi

if [ "${SCREENSHOT_ENABLED}" = "true" ]; then
    printf 'export SCREENSHOT_ENABLED=true\n' >> "${ENV_VARS_FILE}"
fi

if [ -n "${ACCESS_TOKEN}" ]; then
    export HA_ACCESS_TOKEN="${ACCESS_TOKEN}"
    ESPHOME_EXPORTS=$(node /usr/local/bin/discover-esphome.js 2>/dev/null || true)
    if [ -n "${ESPHOME_EXPORTS}" ]; then
        printf '%s\n' "${ESPHOME_EXPORTS}" >> "${ENV_VARS_FILE}"
    fi
fi

# ==============================================================================
# Setup Gemini configuration
# ==============================================================================

GEMINI_USER_DIR="/data/.config/gemini"
mkdir -p "$GEMINI_USER_DIR"
mkdir -p "/data/.cache"
mkdir -p "/data/.local/share"
mkdir -p "/data/.config"

if [ ! -f "$GEMINI_USER_DIR/settings.json" ]; then
    bashio::log.info "Creating default Gemini settings..."
    cat > "$GEMINI_USER_DIR/settings.json" << 'EOF'
{
  "approvalMode": "default",
  "screenReader": false,
  "acp": false,
  "sandbox": false,
  "mcpServers": {}
}
EOF
fi

if [ ! -f "$GEMINI_USER_DIR/trustedFolders.json" ]; then
    bashio::log.info "Creating default trusted folders..."
    cat > "$GEMINI_USER_DIR/trustedFolders.json" << 'EOF'
{
  "/config": "TRUST_FOLDER",
  "/homeassistant": "TRUST_FOLDER",
  "/": "TRUST_FOLDER",
  "/data": "TRUST_FOLDER",
  "/opt": "TRUST_FOLDER"
}
EOF
fi

if [ ! -f "/homeassistant/.geminiignore" ]; then
    cat > "/homeassistant/.geminiignore" << 'EOF'
.storage/
.git/
.gemini*/
backups/
addons/
deps/
local/
share/
tts/
www/
blueprints/
node_modules/
*.db
*.log
EOF
fi

# Copy AGENTS.md template to /homeassistant if missing
if [ ! -f "/homeassistant/AGENTS.md" ]; then
    bashio::log.info "Deploying AGENTS.md template to /homeassistant..."
    cp /opt/ha-mcp-server/AGENTS.md /homeassistant/AGENTS.md
fi

# Generate HA Context if enabled
if [ "${SMART_CONTEXT}" = "true" ]; then
    bashio::log.info "Generating Home Assistant context..."
    /usr/local/bin/ha-context || bashio::log.warning "Failed to generate context"
fi

# Set up MCP server
if [ "$MCP_ENABLED" = "true" ]; then
    bashio::log.info "Configuring MCP server..."
    # The ha-mcp-server is located at /opt/ha-mcp-server/index.js
    jq '.mcpServers["home-assistant"] = {
        "command": "node",
        "args": ["/opt/ha-mcp-server/index.js"],
        "env": {"SUPERVISOR_TOKEN": "'$SUPERVISOR_TOKEN'"}
    }' "$GEMINI_USER_DIR/settings.json" > "$GEMINI_USER_DIR/settings.tmp" && mv "$GEMINI_USER_DIR/settings.tmp" "$GEMINI_USER_DIR/settings.json"
fi

# ==============================================================================
# Install Persistent Packages
# ==============================================================================

# APT Packages
APT_PKG_COUNT=$(bashio::config 'persistent_apt_packages | length')
if [ "${APT_PKG_COUNT}" -gt 0 ]; then
    bashio::log.info "Installing persistent APT packages..."
    APT_PKGS=""
    for i in $(seq 0 $((APT_PKG_COUNT - 1))); do
        PKG=$(bashio::config "persistent_apt_packages[${i}]")
        APT_PKGS="${APT_PKGS} ${PKG}"
    done
    apt-get update >/dev/null
    apt-get install -y --no-install-recommends ${APT_PKGS}
fi

# Python PIP Packages
PIP_PKG_COUNT=$(bashio::config 'persistent_pip_packages | length')
if [ "${PIP_PKG_COUNT}" -gt 0 ]; then
    bashio::log.info "Installing persistent PIP packages..."
    if [ ! -d "/data/.local/share/gemini-venv" ]; then
        python3 -m venv /data/.local/share/gemini-venv
    fi
    printf 'export PATH="/data/.local/share/gemini-venv/bin:$PATH"\n' >> "${ENV_VARS_FILE}"
    PIP_PKGS=""
    for i in $(seq 0 $((PIP_PKG_COUNT - 1))); do
        PKG=$(bashio::config "persistent_pip_packages[${i}]")
        PIP_PKGS="${PIP_PKGS} ${PKG}"
    done
    /data/.local/share/gemini-venv/bin/pip install --no-cache-dir ${PIP_PKGS}
fi

chmod 600 "${ENV_VARS_FILE}"

bashio::log.info "Gemini initialization complete"
